# KaYouWeiXin/8080.conf
[program: KaYouWeiXin_8080]

command = /usr/bin/python3 run_flask.py ;

directory = /home/web/KaYouWeiXin         ; 运行的脚本的目录

autostart = true                                      ;  随supervisor自动启动

autorestart = unexpected                              ; 出错重启

startsecs = 3                         ; 启动后3秒不报错就认为程序启动成功

startretries = 3                      ; 程序失败的重试次数web
